# fs_overlay/etc/services.d/11-inotifywait/run
#!/usr/bin/with-contenv sh

inotifywait -q -m -e close_write --format %e /var/run/domains.d/domains |
while read events; do
  reconfig
done