# fs_overlay/etc/services.d/11-api/run
#!/usr/bin/with-contenv sh

npm start --prefix /src/api/