# fs_overlay/etc/services.d/11-api/run
#!/usr/bin/with-contenv sh

node /src/api/index.js