# fs_overlay/etc/services.d/20-crond/run
#!/usr/bin/with-contenv sh

echo 'Starting crond ...'
cron -f
