# fs_overlay/etc/services.d/20-crond/run
#!/bin/sh

echo 'Starting crond ...'
cron -f
