# root_dir/etc/cron.weekly/renew_certs
#!/usr/bin/ruby
require '/opt/certs_manager/certs_manager'

CertsManager.new.renew
