# fs_overlay/etc/services.d/10-docker-gen/run
#!/bin/sh

echo 'Starting docker-gen ...'
docker-gen -watch -only-exposed -notify reconfig /etc/docker-gen/domains.tmpl /var/run/domains
