# fs_overlay/etc/services.d/00-nginx/run
#!/command/with-contenv sh

nginx -g "daemon off;"
