# fs_overlay/etc/services.d/00-nginx/run
#!/usr/bin/with-contenv sh

nginx -g "daemon off;"
